# script/systemd/skywalker.service
[Unit]
Description=Skywalker
After=network.target
 
[Service]
Type=simple
ExecStart=/usr/local/bin/skywalker -c /etc/skywalker.yml
ExecStop=/usr/local/bin/forctl -c /etc/skywalker.yml quit -y
PrivateTmp=false
 
[Install]
WantedBy=multi-user.target
